# root/etc/s6-overlay/s6-rc.d/plexmediaserver-pid-cleanup/run
#!/command/with-contenv sh

if [ -e /var/lib/plexmediaserver/Library/Application\ Support/Plex\ Media\ Server/plexmediaserver.pid ]; then
  rm -rf /var/lib/plexmediaserver/Library/Application\ Support/Plex\ Media\ Server/plexmediaserver.pid
fi
