# root/etc/s6-overlay/s6-rc.d/plexmediaserver/run
#!/command/with-contenv sh

. /etc/profile.d/plexmediaserver.sh

exec s6-setuidgid appuser /usr/lib/plexmediaserver/Plex\ Media\ Server
