# root/etc/s6-overlay/s6-rc.d/plexmediaserver-fixattrs/run
#!/command/with-contenv sh

chown -R appuser:appuser /transcode
chmod -R u=rwx,g=rx,o-rwx /transcode
