# root/etc/s6-overlay/s6-rc.d/plexmediaserver-transcode/run
#!/command/with-contenv sh

if [ ! -d /transcode ];then
  mkdir -p /transcode
fi
